# ai/ai-server.service
[Unit]
Description=Llama.cpp AI Server for GoStream
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/GoStream/ai
ExecStart=/home/pi/GoStream/ai/llama.cpp/build/bin/llama-server -m /home/pi/GoStream/ai/models/Qwen_Qwen3-0.6B-Q4_K_M.gguf --port 8085 -t 4 -c 512 --host 0.0.0.0 -np 1 --mlock --slot-save-path /tmp/llama-slots
LimitMEMLOCK=infinity
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
